# templates/default/web.conf.erb
pre-start script

bash << "EOF"
  mkdir -p /var/run/unicorn
  touch /var/run/unicorn/master.pid
  chown -R root /var/run/unicorn
EOF

end script

start on starting errata
stop on stopping errata
